# server/deploy/lightsail/proofline-api.service
[Unit]
Description=Proofline Backend API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/proofline/server
EnvironmentFile=/etc/proofline/server.env
ExecStart=/bin/bash -lc 'export NVM_DIR="$HOME/.nvm"; . "$NVM_DIR/nvm.sh"; nvm use 22 >/dev/null; corepack enable >/dev/null 2>&1 || true; pnpm start'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
